# mon_test.service
[Unit]
Description=Monitoring script for test process

[Service]
ExecStart=/usr/local/bin/monitoring_test.sh
